# bedework/supervisord.conf
[supervisord]
nodaemon=true

[program:dirsrv]
command=/bin/bash bw dirstart 

[program:jboss]
command=/bin/bash startjboss -debug

